# DAQ-Multiprocessing-Files-final-pi-tested/unmodified-daq-battery-upload.service
[Unit]
Description=Upload sealed backlog DBs after DAQ stops
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=gholi
Group=gholi
Environment=HOME=/home/gholi
Environment=BLOB_SENDER_DIR=/home/gholi/blob_sender
Environment=PENDING_DIR=/opt/Pilot_Deployment/daq_data/backlog/pending
Environment=UPLOADED_DIR=/opt/Pilot_Deployment/daq_data/backlog/uploaded
Environment=FAILED_DIR=/opt/Pilot_Deployment/daq_data/backlog/failed
Environment=STATE_DIR=/opt/Pilot_Deployment/daq_data/backlog/upload_state
Environment=LOCK_FILE=/opt/Pilot_Deployment/daq_data/backlog/upload.lock
Environment=S3_BUCKET=sonibel-testing
Environment=S3_PREFIX=daq-unmodified-backlog-uploads
Environment=TIME_BUDGET_MIN=30
Environment=MAX_CHUNK_MB=64
Environment=STOP_MARGIN_SECONDS=0
Environment=FOLLOW=1
Environment=CPU_CORE=
WorkingDirectory=/opt/Pilot_Deployment
ExecStart=/opt/Pilot_Deployment/upload_backlog_dbs_to_s3.sh
TimeoutStartSec=32min
StandardOutput=append:/var/log/unmodifiedDaqBatteryUpload.log
StandardError=append:/var/log/unmodifiedDaqBatteryUpload.log

[Install]
WantedBy=multi-user.target
